# scripts/systemd/quantagent-v11-ensemble.service
[Unit]
Description=QuantAgent v11 ensemble launcher (12 fold × 3 seed × Stage 2-5 features)
ConditionPathExists=/home/shanhefu/QuantAgent
OnFailure=quantagent-health-alert@%n.service

[Service]
Type=oneshot
WorkingDirectory=/home/shanhefu/QuantAgent
TimeoutStartSec=24h

# Conda / venv path
Environment=CONDA_PREFIX=/home/shanhefu/QuantAgent/AI_quant_venv

# v11 knobs — override here if needed
Environment=QA_OUTPUT_ROOT=/home/shanhefu/QuantAgent/runtime/models/v7_alpha_v11
Environment=QA_LAKE_ROOT=/home/shanhefu/QuantAgent/runtime/data/v7
Environment=QA_REPLAY_OUT=/home/shanhefu/QuantAgent/runtime/reports/sleeve_replay_v11
Environment=QA_SEEDS=1729,17,42
Environment=QA_SUB_MODELS=main,lowbuy,breakout,limitup_risk

# CUDA selection
Environment=CUDA_VISIBLE_DEVICES=0

ExecStart=/home/shanhefu/QuantAgent/scripts/run_v11_ensemble.sh

StandardOutput=journal
StandardError=journal
SuccessExitStatus=0
